# bin/init.d/pve-storage.target
[Unit]
Description=PVE Storage Target
Wants=remote-fs.target
After=remote-fs.target
After=ceph.service
After=ceph-mon.target
After=ceph-osd.target
After=ceph-mds.target
After=ceph-mgr.target
After=sheepdog.service
After=glusterd.service
After=open-iscsi.service
