# bin/init.d/pve-manager.service
[Unit]
Description=PVE VM Manager
ConditionPathExists=/usr/bin/pvesh
Wants=pveproxy.service
After=pveproxy.service
After=open-iscsi.service

[Service]
ExecStart=/etc/init.d/pve-manager start
ExecStop=/etc/init.d/pve-manager stop
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
