# bin/init.d/pvedaemon.service
[Unit]
Description=PVE API Daemon
ConditionPathExists=/usr/bin/pvedaemon
Wants=corosync.service pve-cluster.service
After=corosync.service pve-cluster.service

[Service]
ExecStart=/usr/bin/pvedaemon start
ExecStop=/usr/bin/pvedaemon stop
ExecReload=/usr/bin/pvedaemon restart
PIDFile=/var/run/pvedaemon.pid
Type=forking

[Install]
WantedBy=multi-user.target
