# bin/init.d/pvedaemon
#!/bin/sh

### BEGIN INIT INFO
# Provides:        pvedaemon
# Required-Start:  $remote_fs $network $syslog pve-cluster cman clvm
# Required-Stop:   $remote_fs $network $syslog pve-cluster cman clvm
# Default-Start:   2 3 4 5
# Default-Stop:    0 1 6
# Short-Description: PVE Daemon
### END INIT INFO

. /lib/lsb/init-functions

PATH=/sbin:/bin:/usr/bin:/usr/sbin
DAEMON=/usr/bin/pvedaemon
NAME=pvedaemon
DESC="PVE Daemon"
PIDFILE=/var/run/pvedaemon.pid

test -f $DAEMON || exit 0

# avoid warnings about uninstalled locales when pvedaemon executes commands
export LC_ALL="C"

case "$1" in
	start)
		log_daemon_msg "Starting $DESC" "$NAME"
		$DAEMON start
		log_end_msg $?
  		;;
	stop)
		log_daemon_msg "Stopping $DESC" "$NAME"
		$DAEMON stop
		log_end_msg $?
  		;;
	restart|reload|force-reload)
		log_daemon_msg "Restarting $DESC" "$NAME"
		$DAEMON restart

 		log_end_msg $?
  		;;
	*)
	        N=/etc/init.d/$NAME
  		echo "Usage: $N {start|stop|restart|reload|force-reload}"
  		exit 1
		;;
esac

exit 0
