# bin/init.d/pveproxy.service
[Unit]
Description=PVE API Proxy Server
ConditionPathExists=/usr/bin/pveproxy
Requires=pve-cluster.service
Wants=pvedaemon.service
After=pve-cluster.service
After=pvedaemon.service

[Service]
ExecStart=/usr/bin/pveproxy start
ExecStop=/usr/bin/pveproxy stop
ExecReload=/usr/bin/pveproxy restart
PIDFile=/var/run/pveproxy/pveproxy.pid
Type=forking

[Install]
WantedBy=multi-user.target
