# bin/init.d/pvesr.service
[Unit]
Description=Proxmox VE replication runner

[Service]
ExecStart=@/usr/bin/pvesr "/usr/bin/pvesr" "run"
