# services/pvesr.service
[Unit]
Description=Proxmox VE replication runner
ConditionPathExists=/usr/bin/pvesr

[Service]
Type=oneshot
ExecStart=/usr/bin/pvesr run --mail 1
