# services/spiceproxy.service
[Unit]
Description=PVE SPICE Proxy Server
ConditionPathExists=/usr/bin/spiceproxy
Wants=pveproxy.service
After=pveproxy.service

[Service]
ExecStart=/usr/bin/spiceproxy start
ExecStop=/usr/bin/spiceproxy stop
ExecReload=/usr/bin/spiceproxy restart
PIDFile=/var/run/pveproxy/spiceproxy.pid
Type=forking

[Install]
WantedBy=multi-user.target
